# supervisord5a.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/sbin/httpd -D FOREGROUND -f '/etc/apache2/httpd.conf'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=0

[program:cron]
command=/usr/sbin/crond -f
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=0
